# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/etc/uci-defaults/luci-glqos
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@glqos[-1]
	add ucitrack glqos
	set ucitrack.@glqos[-1].init=glqos
	commit ucitrack
	delete firewall.glqos
	set firewall.glqos=include
	set firewall.glqos.type=script
	set firewall.glqos.path=/usr/sbin/glqos.sh
	set firewall.glqos.reload=1
	commit firewall
EOF
exit 0
