# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/rootfs-overlay/etc/init.d/gl_fixdomain
#!/bin/sh /etc/rc.common

START=15

start(){
	lan_ip=$(uci get network.lan.ipaddr)
	uci set dhcp.localhost=domain
	uci set dhcp.localhost.name='console.gl-inet.com'
	uci set dhcp.localhost.ip="$lan_ip"
	uci commit dhcp
}


