# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/etc/hotplug.d/tty/20-cimi
#!/bin/sh

. /lib/netifd/netifd-proto.sh

cimi=""
dev=""
case "$DEVICENAME" in
	ttyUSB*)
			/usr/bin/modem-hotplug $DEVPATH $DEVICENAME $ACTION
esac

