# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/rootfs-overlay/etc/uci-defaults/99-disable-upnpd
uci set upnpd.config.enable_upnp=0
uci set upnpd.config.enable_natpmp=0
uci commit upnpd
/etc/init.d/miniupnpd disable
