# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/rootfs-overlay/etc/uci-defaults/dhcp_server
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set dhcp.lan.force=1
	commit dhcp
EOF

exit 0
