# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/gl-fw.list
/etc/config/glfw
/usr/lib/gl/libglfw.so
/usr/lib/libglfw.so
/usr/bin/glfw_unittest
/etc/hotplug.d/iface/70-glfw
/etc/uci-defaults/glfw
/usr/bin/glfw
/etc/init.d/glfw
