# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/gl-tertf.list
/etc/init.d/gl_tertf
/etc/modules.d/98-ter-tertf
/usr/bin/gl_tertf
