# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/gl-wg.list
/etc/init.d/wireguard
/usr/lib/libwgapi.so
/usr/lib/gl/libwgapi.so
/usr/bin/wg_unittest
/etc/config/wireguard
