# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/lighttpd.list
/usr/lib/lighttpd/mod_indexfile.so
/etc/lighttpd/lighttpd.conf
/etc/init.d/lighttpd
/etc/lighttpd/mime.conf
/usr/sbin/lighttpd
/usr/lib/lighttpd/mod_staticfile.so
/usr/lib/lighttpd/mod_dirlisting.so
