# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/smstools3.list
/usr/bin/sendsms
/usr/bin/smsd
/etc/init.d/smstools3
/etc/smsd.conf
