# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/unbound.list
/etc/init.d/unbound
/usr/sbin/unbound
/etc/unbound/unbound_ext.conf
/usr/lib/unbound/rootzone.sh
/usr/lib/unbound/defaults.sh
/usr/lib/unbound/unbound.sh
/usr/lib/unbound/odhcpd.awk
/etc/unbound/unbound_srv.conf
/usr/lib/unbound/dnsmasq.sh
/usr/lib/unbound/iptools.sh
/usr/sbin/unbound-checkconf
/usr/lib/unbound/odhcpd.sh
/etc/config/unbound
/etc/unbound/unbound.conf
/etc/unbound/root.key
/etc/hotplug.d/ntp/25-unbound
