# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/rootfs-overlay/etc/hotplug.d/iface/17-vpn-user
#!/bin/sh

INTERFACES="wan wwan modem tethering wg ovpn"

[ -f "/etc/vpn.user" ] && {
	for interface in ${INTERFACES}; do
		[ "$interface" = "$INTERFACE" ] && {
			env -i ACTION="$ACTION" INTERFACE="$INTERFACE" DEVICE="$DEVICE" \
				/bin/sh /etc/vpn.user
		}
	done
}

exit 0
