# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/gl-base-files.list
/rootfs-overlay/etc/chatscripts/3g.chat
/rootfs-overlay/etc/banner-mt750
/rootfs-overlay/etc/hotplug.d/usb/20-usbmuxd
/rootfs-overlay/etc/vpn.user
/rootfs-overlay/etc/smsd.conf
/rootfs-overlay/etc/hotplug.d/tty/30-3g
/rootfs-overlay/etc/firewall.user
/rootfs-overlay/etc/uci-defaults/99-disable-upnpd
/rootfs-overlay/etc/uci-defaults/101-upgrade_from_2_x_firmware
/rootfs-overlay/etc/sysctl.conf
/rootfs-overlay/etc/lighttpd/conf.d/30-proxy.conf
/rootfs-overlay/etc/lighttpd/mime.conf
/rootfs-overlay/etc/ip-up.d/01-clear-count.sh
/rootfs-overlay/usr/bin/glautoupdater
/rootfs-overlay/etc/gcom/getimsi.gcom
/rootfs-overlay/lib/network/wwan/05c6:9003
/rootfs-overlay/etc/inittab
/rootfs-overlay/usr/bin/switchaction
/rootfs-overlay/lib/netifd/proto/3g.sh
/rootfs-overlay/etc/samba/smb.conf.template
/rootfs-overlay/etc/lighttpd/lighttpd.conf
/rootfs-overlay/etc/init.d/firewall_gl
/rootfs-overlay/lib/network/wwan/05c6:9215
/rootfs-overlay/etc/lighttpd/conf.d/10-port.conf
/rootfs-overlay/etc/hotplug.d/block/10-mount
/rootfs-overlay/etc/init.d/lighttpd
/rootfs-overlay/usr/bin/BTN_1
/rootfs-overlay/etc/uci-defaults/dhcp_server
/rootfs-overlay/usr/bin/button_pressed
/rootfs-overlay/lib/network/wwan/05c6:9090
/rootfs-overlay/usr/bin/make_wisp_relay
/rootfs-overlay/etc/lighttpd/conf.d/30-cgi.conf
/rootfs-overlay/lib/netifd/proto/qmi.sh
/rootfs-overlay/etc/rc.button/BTN_0
/rootfs-overlay/etc/rc.button/BTN_1
/rootfs-overlay/etc/init.d/gl_fixdomain
/rootfs-overlay/usr/bin/wifionoff
/rootfs-overlay/etc/hotplug.d/iface/17-vpn-user
/rootfs-overlay/usr/bin/button_released
/rootfs-overlay/etc/config/glconfig
/rootfs-overlay/usr/bin/apstatus
/rootfs-overlay/etc/ppp/options.pptp
/rootfs-overlay/etc/config/samba
/rootfs-overlay/etc/rc.local
/rootfs-overlay/usr/bin/repeater_onoff
/rootfs-overlay/usr/bin/ovpnchecker.sh
/rootfs-overlay/etc/uci-defaults/99_fix_leds
/rootfs-overlay/etc/hotplug.d/block/50-samba
/rootfs-overlay/etc/rc.button/reset
/rootfs-overlay/etc/init.d/smstools3
/rootfs-overlay/etc/init.d/gl_init
/rootfs-overlay/usr/bin/flashleds
/rootfs-overlay/usr/bin/factorydefault
/rootfs-overlay/usr/bin/sendsms
