# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/gl-monitor.list
/etc/init.d/gl_monitor
/usr/bin/gl_monitor
