# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/gl-qos.list
/etc/cls_u32.ko
/etc/uci-defaults/luci-glqos
/usr/lib/libqosapi.so
/etc/modules.d/99-qos
/usr/lib/gl/libqosapi.so
/etc/xt_u32.ko
/etc/sch_ingress.ko
/etc/init.d/glqos
/etc/hotplug.d/iface/25-glqos
/usr/bin/glqos_unittest
/etc/sch_cbq.ko
/etc/cls_flow.ko
/usr/sbin/glqos.sh
/etc/sch_sfq.ko
/etc/sch_prio.ko
