# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/gl-vpn.list
/usr/lib/gl/libovpnapi.so
/usr/bin/ovpn_check_route
/etc/init.d/startvpn
/usr/lib/libovpnapi.so
/etc/openvpn/update-resolv-conf
