# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/gl-vpn-server.list
/usr/lib/gl/libovpnsrvapi.so
/etc/config/vpn_service
/usr/bin/gen_ovpn
/etc/init.d/vpn-service
/usr/lib/libovpnsrvapi.so
/usr/bin/cert_manager
/etc/openvpn/cert/dh1024.pem
