# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/gl-wg-server.list
/usr/lib/gl/libwgserverapi.so
/etc/init.d/wireguard_server
/etc/config/wireguard_server
/usr/bin/wgserver_unittest
/usr/lib/libwgserverapi.so
