# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/usr/lib/opkg/info/mwan3.list
/etc/mwan3.user
/etc/hotplug.d/iface/15-mwan3
/usr/sbin/mwan3
/lib/mwan3/common.sh
/etc/init.d/mwan3
/etc/hotplug.d/iface/16-mwan3-user
/usr/libexec/rpcd/mwan3
/usr/sbin/mwan3track
/lib/mwan3/mwan3.sh
/etc/uci-defaults/luci-mwan3
/etc/hotplug.d/iface/14-mwan3
/etc/config/mwan3
