# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root-ramips/etc/uci-defaults/glfw
#!/bin/sh

s=firewall.glfw
uci get "$s" >/dev/null && exit 0
uci batch <<-EOF
	set $s=include
	set $s.type=script
	set $s.path=/var/etc/glfw.include
	set $s.reload=1
	commit firewall
EOF

exit 0
