# imagebuilder/3.0/openwrt-imagebuilder-ramips-mt7620/build_dir/target-mipsel_24kc_musl/root.orig-ramips/etc/hotplug.d/iface/70-glfw
#!/bin/sh

[ "$ACTION" = "ifup" -o "$ACTION" = "ifupdate" ] || exit 0

logger -t glfw "Reloading firewall due to $ACTION of $INTERFACE ($DEVICE)"

/etc/init.d/glfw restart
